# runit/RHEL5/upstart-runsvdir.conf
# for runit - manage /etc/runit/2
start on runlevel 2
start on runlevel 3
start on runlevel 4
start on runlevel 5
stop on shutdown
respawn
exec /etc/runit/2
